# KeyboardMapping/rc.local
#!/bin/sh -e 

/usr/bin/loadkeys /etc/keyboard/mymod.map

exit 0
